# wax/sh1mmer_legacy/root/noarch/etc/init/sh1mmer.conf
start on started system-services
exec /usr/sbin/sh1mmer_factory.sh
